# rootfs/etc/services.d/endurian/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Endurian
# Take care of a graceful shutdown of the Endurian service
# ==============================================================================

bashio::log.info "Shutting down Endurian..."
bashio::log.info "Goodbye!"
