# rootfs/etc/services.d/endurian/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Endurian
# Runs the Endurian FastAPI application
# ==============================================================================

bashio::log.info "Starting Endurian application..."

# Get configuration values
DEBUG=$(bashio::config 'debug')
LOG_LEVEL=$(bashio::config 'log_level')
WORKERS=$(bashio::config 'workers')
BIND_HOST=$(bashio::config 'bind_host')
BIND_PORT=$(bashio::config 'bind_port')

# Set logging level based on add-on debug setting
if bashio::config.true 'debug'; then
    export LOG_LEVEL="DEBUG"
    bashio::log.debug "Debug mode enabled"
fi

# Change to backend directory
cd /app/backend || bashio::exit.nok "Failed to change to backend directory"

# Check if backend directory has the expected files
bashio::log.info "Backend directory contents:"
ls -la /app/backend/

# Source environment file
if [[ -f .env ]]; then
    source .env
    bashio::log.debug "Environment file loaded"
else
    bashio::log.warning "Environment file not found, using defaults"
fi

# Check if main.py exists
if [[ ! -f main.py ]]; then
    bashio::exit.nok "main.py not found in /app/backend directory!"
fi

# Check virtual environment
if [[ ! -f /opt/venv/bin/python ]]; then
    bashio::exit.nok "Python virtual environment not found at /opt/venv/bin/python!"
fi

# Set environment variables for the application
export UID=8080
export GID=8080

bashio::log.info "Starting Endurian with the following configuration:"
bashio::log.info "- Host: ${BIND_HOST}"
bashio::log.info "- Port: ${BIND_PORT}"
bashio::log.info "- Workers: ${WORKERS}"
bashio::log.info "- Log Level: ${LOG_LEVEL}"
bashio::log.info "- Debug: ${DEBUG}"

# Add more debugging output first
bashio::log.info "=== DEBUG: Service starting ==="
bashio::log.info "Current user: $(whoami)"
bashio::log.info "Current working directory: $(pwd)"
bashio::log.info "Environment variables:"
env | grep -E "(BIND|LOG_LEVEL|POSTGRES)" || bashio::log.info "No matching env vars"

# Check if required files exist
bashio::log.info "Checking required files..."
if [[ ! -f /opt/venv/bin/python ]]; then
    bashio::log.error "Virtual environment not found at /opt/venv/bin/python"
    ls -la /opt/venv/bin/ || bashio::log.error "Virtual environment directory doesn't exist"
    bashio::exit.nok "Virtual environment missing!"
fi

if [[ ! -f /app/backend/main.py ]]; then
    bashio::log.error "main.py not found at /app/backend/main.py"
    bashio::log.info "Backend directory contents:"
    ls -la /app/backend/ || bashio::log.error "Backend directory doesn't exist"
    bashio::exit.nok "FastAPI main.py missing!"
fi

bashio::log.info "All files exist, starting Endurian..."

# Set environment variables for the application
export UID=8080
export GID=8080

# Start with direct command and explicit output redirection
bashio::log.info "Starting uvicorn directly..."
cd /app/backend

# Use exec without s6-setuidgid first to see if that's the issue
/opt/venv/bin/python -m uvicorn main:app \
    --host "${BIND_HOST}" \
    --port "${BIND_PORT}" \
    --log-level "${LOG_LEVEL,,}"
