# rootfs/etc/services.d/endurian/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Endurian
# Runs the Endurian FastAPI application
# ==============================================================================

# CRITICAL TEST: This should appear in logs if s6 recognizes the service
echo "=== ENDURIAN SERVICE STARTING ===" >&2
echo "This message should appear if s6 recognizes the service" >&2

# Also try bashio logging
bashio::log.info "=== ENDURIAN SERVICE STARTING ==="
bashio::log.info "This message should appear if s6 recognizes the service"

# Basic test - just make sure the service is recognized
sleep 2
bashio::log.info "Service recognition test passed!"

# Get configuration values
DEBUG=$(bashio::config 'debug')
LOG_LEVEL=$(bashio::config 'log_level')
WORKERS=$(bashio::config 'workers')
BIND_HOST=$(bashio::config 'bind_host')
BIND_PORT=$(bashio::config 'bind_port')

bashio::log.info "Configuration loaded:"
bashio::log.info "- Host: ${BIND_HOST}"
bashio::log.info "- Port: ${BIND_PORT}"
bashio::log.info "- Log Level: ${LOG_LEVEL}"

# Check basic paths
bashio::log.info "Checking basic paths..."
bashio::log.info "Current working directory: $(pwd)"
bashio::log.info "Current user: $(whoami)"

# Check if backend directory exists
if [[ -d /app/backend ]]; then
    bashio::log.info "Backend directory exists"
    bashio::log.info "Backend contents:"
    ls -la /app/backend/ | head -10
else
    bashio::log.error "Backend directory missing!"
    bashio::log.info "App directory contents:"
    ls -la /app/ || bashio::log.error "App directory doesn't exist"
fi

# Check virtual environment
if [[ -f /opt/venv/bin/python ]]; then
    bashio::log.info "Virtual environment found"
    /opt/venv/bin/python --version
else
    bashio::log.error "Virtual environment missing!"
    ls -la /opt/ || bashio::log.error "Opt directory doesn't exist"
fi

# Start the actual Endurian application
bashio::log.info "Starting Endurian FastAPI application..."

# Set environment variables for the application from Home Assistant config
export ENDURIAN_UID=8080
export ENDURIAN_GID=8080

# Database configuration
export DB_HOST=$(bashio::config 'postgres_host')
export DB_PORT=$(bashio::config 'postgres_port')
export DB_USER=$(bashio::config 'postgres_user')
export DB_PASSWORD=$(bashio::config 'postgres_password')
export DB_NAME=$(bashio::config 'postgres_db')
export DB_SCHEMA=$(bashio::config 'postgres_schema')

# Application configuration
export SECRET_KEY=$(bashio::config 'secret_key')
export FERNET_KEY=$(bashio::config 'fernet_key')
export APP_TITLE=$(bashio::config 'app_title')
export APP_DESCRIPTION=$(bashio::config 'app_description')
export APP_VERSION=$(bashio::config 'app_version')

# CORS configuration
export CORS_ORIGINS=$(bashio::config 'cors_origins')
export CORS_METHODS=$(bashio::config 'cors_methods')
export CORS_HEADERS=$(bashio::config 'cors_headers')

# Server configuration
export BEHIND_PROXY=$(bashio::config 'behind_proxy')
export WORKERS=$(bashio::config 'workers')
export UPLOAD_FOLDER=$(bashio::config 'upload_folder')
export MAX_CONTENT_LENGTH=$(bashio::config 'max_content_length')

# Logging
export LOG_LEVEL=$(bashio::config 'log_level')
export DEBUG=$(bashio::config 'debug')

bashio::log.info "Environment variables set from Home Assistant config"

# Check if main.py exists
if [[ ! -f /app/backend/main.py ]]; then
    bashio::log.error "main.py not found in /app/backend/"
    bashio::log.info "Available Python files:"
    find /app/backend -name "*.py" | head -10
    bashio::exit.nok "FastAPI main.py missing!"
fi

# Start uvicorn
cd /app/backend
bashio::log.info "Starting uvicorn from /app/backend with main.py"
exec /opt/venv/bin/python -m uvicorn main:app \
    --host "${BIND_HOST}" \
    --port "${BIND_PORT}" \
    --log-level "${LOG_LEVEL,,}"
