# rootfs/etc/services.d/endurian/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Endurian
# Runs the Endurian FastAPI application
# ==============================================================================

# CRITICAL TEST: This should appear in logs if s6 recognizes the service
echo "=== ENDURIAN SERVICE STARTING ===" >&2
echo "This message should appear if s6 recognizes the service" >&2

# Also try bashio logging
bashio::log.info "=== ENDURIAN SERVICE STARTING ==="
bashio::log.info "This message should appear if s6 recognizes the service"

# Basic test - just make sure the service is recognized
sleep 2
bashio::log.info "Service recognition test passed!"

# Get configuration values
DEBUG=$(bashio::config 'debug')
LOG_LEVEL=$(bashio::config 'log_level')
WORKERS=$(bashio::config 'workers')
BIND_HOST=$(bashio::config 'bind_host')
BIND_PORT=$(bashio::config 'bind_port')

bashio::log.info "Configuration loaded:"
bashio::log.info "- Host: ${BIND_HOST}"
bashio::log.info "- Port: ${BIND_PORT}"
bashio::log.info "- Log Level: ${LOG_LEVEL}"

# Check basic paths
bashio::log.info "Checking basic paths..."
bashio::log.info "Current working directory: $(pwd)"
bashio::log.info "Current user: $(whoami)"

# Check if backend directory exists
if [[ -d /app/backend ]]; then
    bashio::log.info "Backend directory exists"
    bashio::log.info "Backend contents:"
    ls -la /app/backend/ | head -10
else
    bashio::log.error "Backend directory missing!"
    bashio::log.info "App directory contents:"
    ls -la /app/ || bashio::log.error "App directory doesn't exist"
fi

# Check virtual environment
if [[ -f /opt/venv/bin/python ]]; then
    bashio::log.info "Virtual environment found"
    /opt/venv/bin/python --version
else
    bashio::log.error "Virtual environment missing!"
    ls -la /opt/ || bashio::log.error "Opt directory doesn't exist"
fi

# For now, just sleep to keep the service running and visible in logs
bashio::log.info "Starting infinite loop for testing..."
while true; do
    bashio::log.info "Endurian service is running... (test mode)"
    sleep 30
done
