# rootfs/etc/services.d/endurian/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Endurian
# Runs the Endurian FastAPI application
# ==============================================================================

bashio::log.info "Starting Endurian application..."

# Get configuration values
DEBUG=$(bashio::config 'debug')
LOG_LEVEL=$(bashio::config 'log_level')
WORKERS=$(bashio::config 'workers')
BIND_HOST=$(bashio::config 'bind_host')
BIND_PORT=$(bashio::config 'bind_port')

# Set logging level based on add-on debug setting
if bashio::config.true 'debug'; then
    export LOG_LEVEL="DEBUG"
    bashio::log.debug "Debug mode enabled"
fi

# Change to backend directory
cd /app/backend || bashio::exit.nok "Failed to change to backend directory"

# Check if backend directory has the expected files
bashio::log.info "Backend directory contents:"
ls -la /app/backend/

# Source environment file
if [[ -f .env ]]; then
    source .env
    bashio::log.debug "Environment file loaded"
else
    bashio::log.warning "Environment file not found, using defaults"
fi

# Check if main.py exists
if [[ ! -f main.py ]]; then
    bashio::exit.nok "main.py not found in /app/backend directory!"
fi

# Check virtual environment
if [[ ! -f /opt/venv/bin/python ]]; then
    bashio::exit.nok "Python virtual environment not found at /opt/venv/bin/python!"
fi

# Set environment variables for the application
export UID=8080
export GID=8080

bashio::log.info "Starting Endurian with the following configuration:"
bashio::log.info "- Host: ${BIND_HOST}"
bashio::log.info "- Port: ${BIND_PORT}"
bashio::log.info "- Workers: ${WORKERS}"
bashio::log.info "- Log Level: ${LOG_LEVEL}"
bashio::log.info "- Debug: ${DEBUG}"

# Switch to endurian user and start the application using virtual environment
exec s6-setuidgid endurian \
    /opt/venv/bin/python -m uvicorn main:app \
        --host "${BIND_HOST}" \
        --port "${BIND_PORT}" \
        --workers "${WORKERS}" \
        --log-level "${LOG_LEVEL,,}" \
        --access-log \
        --loop uvloop \
        --http httptools
