# rootfs/etc/services.d/endurian/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Endurian
# Runs the Endurian FastAPI application
# ==============================================================================

bashio::log.info "Starting Endurian application..."

# Get configuration values
DEBUG=$(bashio::config 'debug')
LOG_LEVEL=$(bashio::config 'log_level')
WORKERS=$(bashio::config 'workers')
BIND_HOST=$(bashio::config 'bind_host')
BIND_PORT=$(bashio::config 'bind_port')

# Set logging level based on add-on debug setting
if bashio::config.true 'debug'; then
    export LOG_LEVEL="DEBUG"
    bashio::log.debug "Debug mode enabled"
fi

# Change to backend directory
cd /app/backend || bashio::exit.nok "Failed to change to backend directory"

# Check if backend directory has the expected files
bashio::log.info "Backend directory contents:"
ls -la /app/backend/

# Source environment file
if [[ -f .env ]]; then
    source .env
    bashio::log.debug "Environment file loaded"
else
    bashio::log.warning "Environment file not found, using defaults"
fi

# Check if main.py exists
if [[ ! -f main.py ]]; then
    bashio::exit.nok "main.py not found in /app/backend directory!"
fi

# Check virtual environment
if [[ ! -f /opt/venv/bin/python ]]; then
    bashio::exit.nok "Python virtual environment not found at /opt/venv/bin/python!"
fi

# Set environment variables for the application
export UID=8080
export GID=8080

bashio::log.info "Starting Endurian with the following configuration:"
bashio::log.info "- Host: ${BIND_HOST}"
bashio::log.info "- Port: ${BIND_PORT}"
bashio::log.info "- Workers: ${WORKERS}"
bashio::log.info "- Log Level: ${LOG_LEVEL}"
bashio::log.info "- Debug: ${DEBUG}"

# Test uvicorn import first
bashio::log.info "Testing uvicorn installation..."
if ! /opt/venv/bin/python -c "import uvicorn; print('uvicorn version:', uvicorn.__version__)"; then
    bashio::exit.nok "uvicorn is not properly installed!"
fi

# Test FastAPI app import
bashio::log.info "Testing FastAPI app import..."
if ! /opt/venv/bin/python -c "from main import app; print('FastAPI app imported successfully')"; then
    bashio::exit.nok "Failed to import FastAPI app from main.py!"
fi

bashio::log.info "All checks passed, starting Endurian application..."

# Export environment variables for the startup script
export BIND_HOST="${BIND_HOST}"
export BIND_PORT="${BIND_PORT}"
export LOG_LEVEL="${LOG_LEVEL,,}"

# Switch to endurian user and start the application using helper script
bashio::log.info "Executing startup script as endurian user..."
exec s6-setuidgid endurian /usr/bin/start-endurian.sh
