# rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Endurain
# Runs nginx reverse proxy for ingress support
# ==============================================================================

bashio::log.info "Starting nginx reverse proxy for ingress support..."

# Create necessary directories
mkdir -p /var/log/nginx /run/nginx

# Start nginx
exec nginx -g "daemon off;" -c /etc/nginx/nginx.conf
